# skel/etc/inittab
::sysinit:/bin/mount -a
::sysinit:/etc/init.d/rcS

::respawn:-/bin/sh
::respawn:/sbin/udhcpd -f /etc/udhcpd.conf
::respawn:/sbin/hostapd /etc/hostapd.conf
/dev/kmsg::respawn:/sbin/adbd

::restart:/sbin/init

# TODO check what follows works and is needed
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a
